# mon/nvidia_exporter.service
[Unit]
After=docker.service

[Service]
ExecStart=/opt/nvidia_exporter

[Install]
WantedBy=default.target
